# launchd/rtk-savings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>rtk-savings</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-a</string>
    <!-- Edit this path if you installed the app somewhere other than /Applications -->
    <string>/Applications/rtk Savings.app</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
